# DevOps/ECommerce-Project/Ansible-RoboShop/roles/user/files/user.service
[Unit]
Description = User Service
[Service]
User=roboshop
Environment=MONGO=true
Environment=REDIS_HOST=redis-dev.learntechnology.cloud
Environment=MONGO_URL="mongodb://mongodb-dev.learntechnology.cloud:27017/users"
ExecStart=/bin/node /app/server.js
SyslogIdentifier=user

[Install]
WantedBy=multi-user.target