# DevOps/ECommerce-Project/Ansible-RoboShop/shipping.service
[Unit]
Description=Shipping Service

[Service]
User=roboshop
Environment=CART_ENDPOINT=cart-dev.learntechnology.cloud:8080
Environment=DB_HOST=mysql-dev.learntechnology.cloud
ExecStart=/bin/java -jar /app/shipping.jar
SyslogIdentifier=shipping

[Install]
WantedBy=multi-user.target
